# 9.2/atomic-contrib/template.service
[Unit]
Description=PostgreSQL 9.2 database server
After=docker.service
BindTo=docker.service

[Service]
ExecStart=/usr/bin/docker start -a TEMPLATE
ExecStop=/usr/bin/docker stop TEMPLATE

[Install]
WantedBy=multi-user.target

